# cron.d/huecron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""
29 16 * * * root export $(cat /tmp/hue.env| xargs);mnrn/my-home-hue/target/release/my_home_hue > /tmp/hue.log 2>&1
